# services/htpcmanager/run
#!/bin/bash
/sbin/setuser abc python /app/htpcmanager/Htpc.py --loglevel error --datadir /config

