# root/etc/services.d/vnstat/run
#!/usr/bin/with-contenv bash
vnstatd -n
