# root/etc/services.d/vnstat/run
#!/usr/bin/with-contenv bash

exec \
	vnstatd -n
